# auto-update.service
[Unit]
Description=Autoupdate

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Su --noconfirm